# et_conf/crontabs/eshtapay
SHELL=/bin/sh
* * * * * /usr/local/eshtapay/watch_firefox.sh 
* * * * * ( sleep 30 && /usr/local/eshtapay/watch_firefox.sh )
*/5 * * * * /usr/local/eshtapay/watch_eshtapay_terminals.sh 
