# et_conf/crontabs/root
SHELL=/bin/sh

*/3 * * * * /home/andrew/repos/et_conf/eshtapay/modem/check_inet.sh
